# etc/ce.conf
[program:ce]
command = python ~/readfree/manage.py celery worker
autostart=true
autorestart=true
redirect_stderr=true
directory=~/readfree/
