# etc/rf.conf
[program:rf]
command = python ~/readfree/manage.py run_gunicorn 0.0.0.0:8001
autostart=true
autorestart=true
redirect_stderr=true
directory=~/readfree/
